# telegram-scraper.service
[Unit]
Description=Telegram Job Scraper
After=network.target
Wants=network.target

[Service]
Type=simple
User=maxbmaapc
WorkingDirectory=/Users/maxbmaapc/GitHub/telegram-job-scraper
Environment=PATH=/Users/maxbmaapc/GitHub/telegram-job-scraper/venv/bin
ExecStart=/Users/maxbmaapc/GitHub/telegram-job-scraper/venv/bin/python run.py --mode continuous
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
